# scripts/aslam-update.timer
[Unit]
Description=Check for Aslam updates every 5 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
