# scripts/aslam-update.service
[Unit]
Description=Aslam auto-update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/home/aslam/scripts/deploy.sh main
Environment=HOME=/home/aslam
Environment=PATH=/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
