# scripts/aslam.service
[Unit]
Description=Aslam AI Assistant
After=network.target

[Service]
Type=simple
User=aslam
Group=aslam
WorkingDirectory=/home/aslam
ExecStart=/home/aslam/aslam
Restart=always
RestartSec=5
Environment=HOME=/home/aslam
EnvironmentFile=/home/aslam/.env

[Install]
WantedBy=multi-user.target
