# scripts/nasir.service
[Unit]
Description=Nasir Knowledge Base
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/nasir
ExecStart=/home/exedev/nasir/nasir
Restart=always
RestartSec=5
Environment=HOME=/home/exedev

[Install]
WantedBy=multi-user.target
